# runit/bot/run
#!/usr/bin/env bash

cd /winey
exec python3 -m winey.bot
