# runit/webapp/run
#!/usr/bin/env bash

cd /wine_log
exec python3 -m wine_log.webapp

